# build/supervisor.ini
[inet_http_server]
port = :9002

[program:socks_server]
directory = /usr/local/services/
command = ./socks_server
stderr_logfile=/usr/local/services/logs/server_stderr.log,/dev/stderr
stdout_logfile=/usr/local/services/logs/server_stdout.log,/dev/stdout
autorestart = true